# cookbooks/arcgis-enterprise/templates/default/arcgisdatastore.service.erb
# ------------------------------------------------------------------
# ArcGIS Data Store systemd unit file
# ------------------------------------------------------------------
#
# Configure ArcGIS Data Store to be started at boot on RHEL 7.x by
# following these instructions:
#
# 1.) Switch to the root user.
#
# 2.) Copy this file to /etc/systemd/system
#
# 3.) Enable the service to start at boot:
#
#     # systemctl enable arcgisdatastore.service
#
# 4.) Verify systemd service is setup correctly:
#
#    # systemctl stop arcgisdatastore.service
#    # systemctl start arcgisdatastore.service
#    # systemctl status arcgisdatastore.service
#
# 5.) Reboot the system and verify that Server restarts properly.
#
# ------------------------------------------------------------------

[Unit]
Description=ArcGIS Data Store Service
After=network.target

[Service]
Type=forking
User=<%=@agsuser%>
GuessMainPID=false
# Configuration for spatiotemporal big data store
LimitNOFILE=65536
LimitNPROC=25059
LimitAS=infinity
LimitFSIZE=infinity
ExecStart=<%=@datastorehome%>/startdatastore.sh
ExecStop=<%=@datastorehome%>/stopdatastore.sh

[Install]
WantedBy=multi-user.target

