# cookbooks/arcgis-enterprise/templates/default/arcgisserver.service.erb
# ------------------------------------------------------------------
# ArcGIS Server systemd unit file
# ------------------------------------------------------------------
#
# Configure ArcGIS Server to be started at boot on RHEL 7.x by
# following these instructions:
#
# 1.) Switch to the root user.
#
# 2.) Copy this file to /etc/systemd/system
#
# 3.) Enable the service to start at boot:
#
#     # systemctl enable arcgisserver.service
#
# 4.) Verify systemd service is setup correctly:
#
#    # systemctl stop arcgisserver.service
#    # systemctl start arcgisserver.service
#    # systemctl status arcgisserver.service
#
# 5.) Reboot the system and verify that Server restarts properly.
#
# ------------------------------------------------------------------

[Unit]
Description=ArcGIS Server Service
After=network.target

[Service]
Type=forking
User=<%=@agsuser%>
GuessMainPID=false
ExecStart=<%=@agshome%>/startserver.sh
ExecStop=<%=@agshome%>/stopserver.sh

[Install]
WantedBy=multi-user.target

